# deploy/supervisor-wordmaster.conf
; WordMaster Supervisor 配置文件
; 复制到 /etc/supervisor/conf.d/wordmaster.conf

[program:wordmaster-backend]
directory=/opt/wordmaster/backend
; 使用 gunicorn 生产服务器，绑定本地地址
command=/opt/wordmaster/backend/venv/bin/gunicorn app.main:app -w 4 -k uvicorn.workers.UvicornWorker --bind 127.0.0.1:8000 --access-logfile /var/log/wordmaster/backend.access.log --error-logfile /var/log/wordmaster/backend.error.log
autostart=true
autorestart=true
user=www-data
stderr_logfile=/var/log/wordmaster/backend.err.log
stdout_logfile=/var/log/wordmaster/backend.out.log
; 生产环境变量
environment=PYTHONPATH="/opt/wordmaster/backend",ENV="production",SECRET_KEY="your-secret-key-change-this"
; 进程管理
stopasgroup=true
killasgroup=true
; 启动等待时间
startsecs=5
; 重启间隔
startretries=3

; 注意：前端不需要 Supervisor 管理！
; 前端使用 Nginx 直接托管静态文件（dist/ 目录）
; 不要运行 npm run dev 或 npm run preview
